# hackyourgrade/conf/apiecorp.service
[Unit]
Description=ECorp REST API Consumer (fetch unread message count)

[Service]
Type=Simple
ExecStart=/root/apiecorp.sh

